# crontab
MAILTO=""

40 * * * * /home/jbrandt/joelbrandt.org/howsiandoing/ian.sh
